# scripts/pi/usb-status-bridge.service
[Unit]
Description=USB Mobile Status Bridge
After=network.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
Environment=HOME=/home/__APP_USER__
Environment=APP_USER=__APP_USER__
Environment=APP_REPO_DIR=/home/__APP_USER__/DuckyDu
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/__APP_USER__/.Xauthority
ExecStart=/usr/local/bin/usb-status-bridge
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
